# files/general/etc/services.d/apache2/run
#!/usr/bin/with-contenv sh
exec httpd -D FOREGROUND
